# deploy/ferienplanung-backend.service
[Unit]
Description=Ferienplanung Backend (Node/Express)
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/root/Ferienplanner/server
EnvironmentFile=/etc/ferienplaner/ferienplaner.env
ExecStart=/usr/bin/node /root/Ferienplanner/server/server.js
Restart=on-failure
RestartSec=2
KillMode=mixed
TimeoutStopSec=10

StateDirectory=ferienplaner

[Install]
WantedBy=multi-user.target
